# debian/araneum.postrm
#!/bin/sh

#DEBHELPER#
update-rc.d -f php-fcgi remove >/dev/null || exit $?
